# bin/meercat
#!/bin/bash
#
# meercat
#
# chkconfig: 345 95 25
# description: meercat solr script
#
if [ "${SOLR_HOME}" == "" -o "${SOLR_DIR}" == "" ]; then 
    [ -f /etc/sysconfig/meercat ] &&  . /etc/sysconfig/meercat
fi

if [ "${SOLR_STOP_PORT}" == "" ]; then
    SOLR_STOP_PORT=8079
fi

if [ "${SOLR_START_PORT}" == "" ]; then
    SOLR_START_PORT=8983
fi

if [ "${SOLR_MEMORY}" == "" ]; then
    SOLR_MEMORY="2048m"
fi

if [ "${SOLR_LOG}" == "" ]; then
    SOLR_LOG=/dev/null
fi

echo "Solr Home: ${SOLR_HOME}"
echo "Solr Dir: ${SOLR_DIR}"
echo "Solr Memory: ${SOLR_MEMORY}"
echo "Solr Port: ${SOLR_START_PORT} + ${SOLR_STOP_PORT}"
echo "Solr Log: ${SOLR_LOG}"

export JETTY_ARGS=/opt/catmandu-meercat-solr/etc/jetty-plus.xml

JAVA_OPTIONS="-d64 -server -Xmx${SOLR_MEMORY} -XX:+AggressiveOpts -XX:NewRatio=5 -Dmeercat -Djetty.port=${SOLR_START_PORT} -DSTOP.PORT=${SOLR_STOP_PORT} -DSTOP.KEY=stopkey -Dsolr.solr.home=${SOLR_HOME} -jar start.jar"
JAVA="/usr/bin/java"

case $1 in
    debug) 
        echo "Starting Solr"
        cd ${SOLR_DIR}
        $JAVA $JAVA_OPTIONS 
        ;;
    start)
        echo "Starting Solr"
        cd ${SOLR_DIR}
        $JAVA "-Djava.util.logging.config.file=etc/logging.properties" $JAVA_OPTIONS 2> ${SOLR_LOG} &
        ;;
    stop)
        echo "Stopping Solr"
        cd ${SOLR_DIR}
        $JAVA $JAVA_OPTIONS --stop
        ;;
    log)
        tail -f ${SOLR_LOG} 
        ;;
    status)
        ps -ef | grep "[m]eercat"
        ;;
    restart)
        $0 stop
        sleep 1
        $0 start
        ;;
    *)
        echo "Usage: $0 {debug|start|stop|restart|status|log}" >&2
        exit 1
        ;;
esac
